# profile.d/00_init
STFILESIZE=10000
HISTSIZE=10000

PATH="$HOME/bin:$PATH"

export CLICOLOR="true"

